# .config/systemd/user/git-keepass-pull.path
[Path]
PathModified=/tmp/keepassxc-mateusz.lock

[Install]
WantedBy=default.target

