# .config/systemd/user/git-keepass-pull.timer
[Timer]
OnBootSec=30sec
OnUnitActiveSec=30min

[Install]
WantedBy=default.target

