# .config/systemd/user/git-keepass-push.path
[Path]
PathChanged=/home/mateusz/.local/share/git-keepass/pass.kdbx

[Install]
WantedBy=default.target

