# .config/systemd/user/git-keepass-push.path
[Path]
PathChanged=/home/mateusz/data/keepass/pass.kdbx

[Install]
WantedBy=default.target

